# tqqq_bot_v5/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0

[program:ibgateway]
command=xvfb-run --auto-servernum /opt/ibc/gatewaystart.sh 9999 -inline --tws-path=/root/Jts --tws-settings-path=/root/Jts
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:botpy]
command=bash -c "python3 /app/gateway/wait_for_gateway.py --port $IBKR_PORT && python3 /app/main.py"
autostart=true
autorestart=true
startsecs=30
depends_on=ibgateway
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
